# build/linux/fim-client.service
[Unit]
Description=FIM Admin Daemon - Cross-platform File Integrity Monitoring
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/opt/fim-client/fim_admin
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target