# src/linux/build/fim-daemon_1.0.0_amd64/etc/systemd/system/fim-daemon.service
[Unit]
Description=File Integrity Monitor Daemon
After=network.target

[Service]
Type=simple
User=fim-daemon
Group=fim-daemon
ExecStart=/usr/bin/python3 /usr/local/bin/fim-daemon
WorkingDirectory=/var/lib/fim-daemon
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target