# src/linux/packaging/DEBIAN/prerm
#!/bin/bash
set -e

echo "Stopping FIM Daemon service..."
systemctl stop fim-daemon || true
systemctl disable fim-daemon || true